# .local/backup/trackpoint-sensitivity.service
[Unit]
Description=Script

[Service]
ExecStart=/home/chris/.local/bin/trackpoint-adjust

[Install]
WantedBy=multi-user.target 
